# etc/supervisor/ludovicopratesi.conf
[program:ludovicopratesi]
directory=/srv/apps/ludovicopratesi
command=/bin/sh -c 'set -a; . /srv/apps/ludovicopratesi/.env.production; exec /srv/venvs/ludovicopratesi/bin/gunicorn www.wsgi:application --bind 127.0.0.1:8002 --workers 3 --timeout 60'
user=deploy
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/logs/ludovicopratesi-gunicorn.log
environment=PYTHONUNBUFFERED="1",DJANGO_SETTINGS_MODULE="www.settings.prod"
